# gemini-web.service
[Unit]
Description=Gemini Web API Wrapper
Documentation=https://github.com/iiroak/gemini-web-api-wrapper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/opt/gemini-web/.venv/bin/gemini-web check
ExecStart=/opt/gemini-web/.venv/bin/gemini-web serve
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt/gemini-web

# Run as a dedicated user (create with: sudo useradd -r -s /bin/false gemini-web)
# User=gemini-web
# Group=gemini-web

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/.gemini-web
PrivateTmp=true

# Environment
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
