# gemini-web.service
[Unit]
Description=Gemini Web API Wrapper
Documentation=https://github.com/iiroak/gemini-web-api-wrapper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=GEMINI_WEB_HOME=/opt/gemini-web/data
ExecStartPre=/opt/gemini-web/.venv/bin/gemini-web check
ExecStart=/opt/gemini-web/.venv/bin/gemini-web serve
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt/gemini-web

# User and Group are set by install.sh based on who runs the installer
# User=__USER__
# Group=__GROUP__

[Install]
WantedBy=multi-user.target
